# test/packaging/armel/DEBIAN-FW/postrm
#!/bin/sh

update-rc.d -f hipfw remove 



