# test/packaging/armel/DEBIAN-FW/preinst
#!/bin/sh
invoke-rc.d --force --quiet hipfw stop 
update-rc.d -f hipfw remove 

 
