# test/packaging/armel/DEBIAN-FW/prerm
#!/bin/sh

invoke-rc.d --force --quiet hipfw stop 


