# test/packaging/armel/DEBIAN-dnsproxy/postrm
#!/bin/sh

update-rc.d -f hipdnsproxy remove 

exit 0



