# test/packaging/armel/DEBIAN-dnsproxy/prerm
#!/bin/sh

invoke-rc.d --force --quiet hipdnsproxy stop 

exit 0


