# tools/daemontools/admin/daemontools-0.76/package/boot.upstart
# daemontools
#
# This service maintains a daemontools' svscanboot

start on stopped rc2
start on stopped rc3
start on stopped rc4
start on stopped rc5

stop on runlevel 0
stop on runlevel 1
stop on runlevel 6

respawn
exec /usr/sbin/svscanboot
