# test/packaging/rh-init.d-dnsproxy
#!/bin/bash
#
# Init file for hip dns proxy
#
# description: HIP DNS proxy
#
# processname: dnsproxy.py
# pidfile: /var/run/dnsproxy.pid
# chkconfig: 2345 22 78

# source function library
. /etc/rc.d/init.d/functions

RETVAL=0
DNSPROXY="/usr/sbin/hipdnsproxy"

PID_FILE=/var/run/hipdnsproxy.pid
LOCK_FILE=/var/lock/subsys/hipdnsproxy
OPTIONS="-bk"

export PATH=/usr/sbin:$PATH

# pull in sysconfig settings
[ -f /etc/sysconfig/hipdnsproxy ] && . /etc/sysconfig/hipdnsproxy

start()
{
	echo -n $"Starting $BASENAME:"
	$DNSPROXY $OPTIONS && success || failure
	RETVAL=$?
	[ "$RETVAL" = 0 ] && touch $LOCK_FILE
	echo
}

stop()
{
	echo -n $"Stopping $BASENAME:"
	kill `head -1 $PID_FILE`
	RETVAL=$?
	[ "$RETVAL" = 0 ] && rm -f $LOCK_FILE
	echo
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		sleep 3
		start
		;;
	status)
		status $DNSPROXY
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
		RETVAL=1
esac
exit $RETVAL
