# test/packaging/armel/DEBIAN/postrm
#!/bin/sh

update-rc.d -f hipd remove 


