# test/packaging/armel/DEBIAN/preinst
#!/bin/sh
invoke-rc.d --force --quiet hipd stop 
update-rc.d -f hipd remove 
 
