# test/packaging/armel/DEBIAN/prerm
#!/bin/sh

invoke-rc.d --force --quiet hipd stop 


